# archiso/airootfs/etc/systemd/system/flynn-night.service
[Unit]
Description=Flynn OS evening mode (warm colors + DND)
[Service]
Type=oneshot
ExecStart=/usr/local/bin/flynn-night night
