# archiso/airootfs/etc/systemd/system/flynn-daemon.service
[Unit]
Description=Flynn OS Daemon — REST API + Pi bridge + ANTIGRAVITY IPC
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/flynn/daemon/flynn_daemon.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
