# archiso/airootfs/etc/systemd/system/flynn-day.service
[Unit]
Description=Flynn OS day mode (colors + notifications back to normal)
[Service]
Type=oneshot
ExecStart=/usr/local/bin/flynn-night day
