# archiso/airootfs/etc/systemd/system/flynn-night.timer
[Unit]
Description=Flynn OS evening mode at 21:00
[Timer]
OnCalendar=*-*-* 21:00
Persistent=false
[Install]
WantedBy=timers.target
