# archiso/airootfs/etc/systemd/system/flynn-update.timer
[Unit]
Description=Flynn OS Nightly Update Timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
