# archiso/airootfs/root/.config/autostart/pipewire.desktop
[Desktop Entry]
Type=Application
Name=PipeWire
Exec=pipewire
